# ironmime.service
[Unit]
Description=IronMime daemon

[Service]
Type=simple
Restart=always
RestartSec=3
ExecStart=/usr/bin/ironmime
User=
Group=input

[Install]
WantedBy=default.target